# dist/copyzen.service
[Unit]
Description=copyzen clipboard recorder (wl-paste watcher)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# wl-paste options MUST precede --watch (everything after --watch is the command).
ExecStart=/usr/bin/env wl-paste --no-newline --watch copyzen store
Restart=on-failure

[Install]
WantedBy=graphical-session.target
